# etc/init.d/supervisor-openrc
#!/sbin/openrc-run

description="Start/stop supervisor"
pidfile="/var/run/supervisord.pid"
command="/usr/bin/supervisord"
command_args="-c /etc/supervisor/supervisord.conf"
retry="TERM/30/KILL/5"

depend() {
    need mountall
    after bootmisc
}

start_pre() {
    checkpath -f -m 0644 -o root:root "$pidfile"
    # Include supervisor defaults if available
    if [ -f /etc/default/supervisor ]; then
        . /etc/default/supervisor
    fi
}

start() {
    ebegin "$description"
    start-stop-daemon --start --quiet --background --pidfile "$pidfile" \
        --exec "$command" -- $command_args
    eend $?
}

stop() {
    ebegin "$description"
    start-stop-daemon --stop --quiet --oknodo --pidfile "$pidfile"
    eend $?
}

restart() {
    ebegin "Restarting $description"
    start-stop-daemon --stop --retry "$retry" --quiet --oknodo --pidfile "$pidfile"
    start-stop-daemon --start --quiet --background --pidfile "$pidfile" \
        --exec "$command" -- $command_args
    eend $?
}

status() {
    status_of_proc -p "$pidfile" "$command" && return 0 || return $?
}
