# debian/postinst
#!/usr/bin/env bash

# -- Add service to runlevel.

rc-update add supervisor-openrc async

# -- Update services list.

rc-update -u
